# configs/power/systemd-user/g14-power-startup-eco.service
[Unit]
Description=G14 power startup default (Eco)
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/bin/bash -lc 'for i in $(seq 1 30); do powerprofilesctl set power-saver && break; sleep 2; done; %h/.local/bin/g14-power-mode.sh apply --logout-on-pending no || true'

[Install]
WantedBy=default.target
